# miyo.service
[Unit]
Description=miyo

[Service]
WorkingDirectory=/opt/miyo
ExecStart=/opt/miyo/miyo

[Install]
WantedBy=default.target
